# backend/supervisor.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:yii-queue-worker]
command=php /app/yii queue/listen --verbose=1 --color=0
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/app/runtime/queue.log
